# crontab
* * * * * php /var/www/html/app/artisan schedule:run >> /dev/null 2>&1
